# crontab
* * * * * ubuntu /bin/bash /home/ubuntu/apps/fy17report/repository/run_on_server.sh fab $DEPLOYMENT_TARGET cron_jobs.test >/var/log/fy17report/crontab.log 2>&1
